# service/new_horizon_gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/new_horizon_gunicorn.sock

[Install]
WantedBy=sockets.target